# docs/installer/services/linux/giljo-mcp.service
[Unit]
Description=GiljoAI MCP Application Service
Documentation=https://github.com/GiljoAI/giljo-mcp
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service
Requires=network.target

[Service]
Type=simple
User=giljo
Group=giljo
WorkingDirectory=/home/giljo/.giljo_mcp

# Main application command
ExecStart=/usr/bin/python3 -m src.giljo_mcp.server

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Environment variables
Environment="GILJO_CONFIG_DIR=/home/giljo/.giljo_mcp/config"
Environment="GILJO_DATA_DIR=/home/giljo/.giljo_mcp/data"
Environment="GILJO_LOG_DIR=/home/giljo/.giljo_mcp/logs"
Environment="PYTHONPATH=/home/giljo/.giljo_mcp"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=giljo-mcp

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/giljo/.giljo_mcp

# Process limits
LimitNOFILE=1024
LimitNPROC=512

# Memory and CPU limits (optional)
MemoryMax=1G
CPUQuota=50%

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

# PID file
PIDFile=/home/giljo/.giljo_mcp/run/giljo-app.pid

[Install]
WantedBy=multi-user.target

# 
# Additional service files should be created:
# 
# 1. giljo-postgresql.service:
#    [Unit]
#    Description=PostgreSQL database for GiljoAI MCP
#    After=network.target
#    
#    [Service]
#    Type=notify
#    User=postgres
#    ExecStart=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/giljo
#    ExecReload=/bin/kill -HUP $MAINPID
#    KillMode=mixed
#    KillSignal=SIGINT
#    TimeoutSec=120
#    
#    [Install]
#    WantedBy=multi-user.target
# 
# 2. giljo-redis.service:
#    [Unit]
#    Description=Redis cache for GiljoAI MCP
#    After=network.target
#    
#    [Service]
#    Type=notify
#    User=redis
#    Group=redis
#    ExecStart=/usr/bin/redis-server /etc/redis/giljo-redis.conf
#    ExecStop=/bin/kill -s QUIT $MAINPID
#    TimeoutStopSec=0
#    Restart=always
#    
#    [Install]
#    WantedBy=multi-user.target
# 
# 3. giljo-worker.service:
#    [Unit]
#    Description=GiljoAI MCP Background Worker
#    After=network.target giljo-mcp.service
#    Wants=giljo-mcp.service
#    
#    [Service]
#    Type=simple
#    User=giljo
#    Group=giljo
#    WorkingDirectory=/home/giljo/.giljo_mcp
#    ExecStart=/usr/bin/python3 -m src.giljo_mcp.worker
#    Restart=always
#    RestartSec=15
#    Environment="GILJO_CONFIG_DIR=/home/giljo/.giljo_mcp/config"
#    Environment="GILJO_DATA_DIR=/home/giljo/.giljo_mcp/data"  
#    Environment="GILJO_LOG_DIR=/home/giljo/.giljo_mcp/logs"
#    
#    [Install]
#    WantedBy=multi-user.target
#
# Installation Instructions:
# 1. Copy service files to /etc/systemd/system/
# 2. Reload systemd: sudo systemctl daemon-reload
# 3. Enable services: sudo systemctl enable giljo-mcp.service
# 4. Start services: sudo systemctl start giljo-mcp.service
# 5. Check status: sudo systemctl status giljo-mcp.service
# 
# For user services (non-root):
# 1. Copy to ~/.config/systemd/user/
# 2. systemctl --user daemon-reload
# 3. systemctl --user enable giljo-mcp.service
# 4. systemctl --user start giljo-mcp.service