# docs/installer/services/macos/com.giljo.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Main GiljoAI MCP Application Service -->
    <key>Label</key>
    <string>com.giljo.mcp.application</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.giljo_mcp.server</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/USER_HOME/.giljo_mcp</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>GILJO_CONFIG_DIR</key>
        <string>/Users/USER_HOME/.giljo_mcp/config</string>
        <key>GILJO_DATA_DIR</key>
        <string>/Users/USER_HOME/.giljo_mcp/data</string>
        <key>GILJO_LOG_DIR</key>
        <string>/Users/USER_HOME/.giljo_mcp/logs</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/USER_HOME/.giljo_mcp/logs/application.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/USER_HOME/.giljo_mcp/logs/application.err</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Service Dependencies -->
    <key>LaunchOnlyOnce</key>
    <false/>
    
    <!-- Resource Limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
        <key>NumberOfProcesses</key>
        <integer>100</integer>
    </dict>
    
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
        <key>NumberOfProcesses</key>
        <integer>200</integer>
    </dict>
    
    <!-- Network settings for the service -->
    <key>NetworkState</key>
    <true/>
    
    <!-- Nice value (process priority) -->
    <key>Nice</key>
    <integer>0</integer>
    
    <!-- Exit timeout -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
    
    <!-- 
    NOTE: This is a template plist file for the main GiljoAI MCP application.
    
    For PostgreSQL service, create a separate plist with:
    - Label: com.giljo.mcp.postgresql
    - ProgramArguments: ["/usr/local/bin/postgres", "-D", "/Users/USER_HOME/.giljo_mcp/data/postgresql"]
    
    For Redis service, create a separate plist with:
    - Label: com.giljo.mcp.redis  
    - ProgramArguments: ["/usr/local/bin/redis-server", "/Users/USER_HOME/.giljo_mcp/config/redis.conf"]
    
    For Worker service, create a separate plist with:
    - Label: com.giljo.mcp.worker
    - ProgramArguments: ["/usr/bin/python3", "-m", "src.giljo_mcp.worker"]
    - RunAtLoad: false (optional service)
    
    Installation Instructions:
    1. Replace USER_HOME with actual user home directory
    2. Copy to ~/Library/LaunchAgents/
    3. Load with: launchctl load ~/Library/LaunchAgents/com.giljo.mcp.plist
    4. Start with: launchctl start com.giljo.mcp.application
    
    Uninstallation:
    1. launchctl stop com.giljo.mcp.application
    2. launchctl unload ~/Library/LaunchAgents/com.giljo.mcp.plist
    3. rm ~/Library/LaunchAgents/com.giljo.mcp.plist
    -->
</dict>
</plist>